# crontab
# Run the script every hour by default
0 * * * * cd /app && python download_frigate.py >> /var/log/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file. 